# Data Audit/Systemd/systemd-lob-quality-daily.service
[Unit]
Description=Daily LOB quality reports
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=lob
Group=lob
WorkingDirectory=/opt/lob_system
ExecStart=/opt/lob_system/venv/bin/python /opt/lob_system/scripts/run_daily_quality_reports_beta.py --root /opt/lob_system --skip-today
