# Data Audit/Systemd/systemd-lob-capture-binance-backup.service
[Unit]
Description=LOB capture Binance BTCUSDT backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lob
Group=lob
WorkingDirectory=/opt/lob_system
ExecStart=/opt/lob_system/venv/bin/python /opt/lob_system/scripts/live_lob_capture.py --exchange binance --storage-exchange binance_backup --symbol BTCUSDT --rotate-hours 12 --root /opt/lob_system
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
