# Data Audit/Systemd/systemd-lob-capture-binance.service
[Unit]
Description=LOB raw capture - Binance BTCUSDT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lob
Group=lob
WorkingDirectory=/opt/lob_system
ExecStart=/opt/lob_system/venv/bin/python /opt/lob_system/scripts/live_lob_capture.py --exchange binance --symbol BTCUSDT --rotate-hours 12 --root /opt/lob_system
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=120
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
